# systemd/microxrce-dds-agent.service
[Unit]
Description=Micro XRCE-DDS Agent (serial on ttyTHS1 for PX4/ROS 2)
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/MicroXRCEAgent serial --dev /dev/ttyTHS1 -b 921600
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
